# build/lib/pyedurov2/system/pyedurov2.service
[Unit]
Description=Run Edu-rover server

[Service]
ExecStart=pyedurov2

SystemCallArchitectures=native
MemoryDenyWriteExecute=true
NoNewPrivileges=true

[Install]
WantedBy=default.target